# ros2_hc_drv/sensomative_ros/hc-sensomative-driver.service
[Unit]
Description=sensomative driver
After=network.target

[Service]
User=scai
Environment=/home/scai/ros4hc.env
ExecStart=/bin/bash -c 'source /opt/ros/current/setup.bash && source install/setup.bash && ros2 launch sensomative_ros sensomative_node.launch.py'
Restart=always

[Install]
WantedBy=multi-user.target
